# Resources/tech.reidar.vifty.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>tech.reidar.vifty.daemon</string>
	<key>BundleProgram</key>
	<string>Contents/MacOS/ViftyDaemon</string>
	<key>MachServices</key>
	<dict>
		<key>tech.reidar.vifty.daemon</key>
		<true/>
	</dict>
	<key>RunAtLoad</key>
	<false/>
	<key>LowPriorityIO</key>
	<true/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>Nice</key>
	<integer>10</integer>
	<key>Umask</key>
	<integer>63</integer>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>StandardOutPath</key>
	<string>/var/log/tech.reidar.vifty.daemon.out.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/tech.reidar.vifty.daemon.err.log</string>
</dict>
</plist>
